# skills/im-bot-connector/templates/supervisor-imbot.conf
[program:hermes-imbot-default]
command=python3 /root/.local/bin/hermes-imbot-listener
autorestart=true
autostart=true
startsecs=5
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/hermes-imbot-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=IMBOT_BACKEND="hermes",IMBOT_TIMEOUT="180"

[program:hermes-imbot-yiman]
command=python3 /root/.hermes/profiles/yiman/home/.local/bin/hermes-imbot-listener
autorestart=true
autostart=true
startsecs=5
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/hermes-imbot-yiman.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/root/.hermes/profiles/yiman",HERMES_HOME="/root/.hermes/profiles/yiman",IMBOT_BACKEND="hermes",IMBOT_STALE_AFTER="120"
